# deploy/secunit-io.service
# Example systemd unit — copy to /etc/systemd/system/secunit-io.service, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now secunit-io
#
# App root (Astro dist): /opt/secunit.io/web
# @astrojs/node emits a Node server — use Node.js LTS for production (see deploy/README.md).
[Unit]
Description=Secunit website (Astro SSR)
After=network.target

[Service]
Type=simple
User=secunit
Group=secunit
WorkingDirectory=/opt/secunit.io/web
Environment=HOST=127.0.0.1
Environment=PORT=4321
Environment=NODE_ENV=production
ExecStart=/opt/secunit.io/.nvm/versions/node/v26.0.0/bin/node /opt/secunit.io/web/server/entry.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
